# uninstaller/etc/udev/rules.d/98-autorunHook.rules
KERNEL=="loop0", RUN+="/root/autorun.sh"
